# setup/configs/shadowplug-service
#!/bin/sh /etc/rc.common
START=99

start() {
    cd /opt/shadowplug/web
    if [ -f app.so ]; then
        python3 -B -c "import app; app.app.run(host='0.0.0.0', port=80)" >> /var/log/shadowplug.log 2>&1 &
    else
        python3 -B app.py >> /var/log/shadowplug.log 2>&1 &
    fi
    /usr/bin/usb-watchdog >> /var/log/shadowplug.log 2>&1 &
    logger "shadowplug: started"
}

stop() {
    killall -9 python3 2>/dev/null
    killall usb-watchdog 2>/dev/null
    logger "shadowplug: stopped"
}
